# config/upstart.conf
description "Cook County Jail API 2.0-dev"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5
exec /home/ubuntu/website/2.0/websites/active/gunicorn.sh
